# bin/minifi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.apache.nifi.minifi</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/minifi-cpp/bin/minifi</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MINIFI_HOME</key>
        <string>/opt/minifi-cpp</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/var/log/minifi.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/minifi.err</string>
</dict>
</plist>
