# bin/minifi.service.in
[Unit]
Description=Starts and stops the MiNiFi C++ Agent
After=network.target

[Service]
Type=simple
Environment=MINIFI_HOME=@MINIFI_HOME_ENV_VALUE@
ExecStart=@MINIFI_EXEC_PATH@
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=20
RestartForceExitStatus=3

[Install]
WantedBy=multi-user.target
